# dot_config/systemd/user/backup.timer
[Unit]
Description=Backup at least every Day

[Timer]
OnBootSec=3min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
