# dot_config/systemd/user/morning.service
[Unit]
Description=Worship
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/skripte/
ExecStart=/root/skripte/morning.sh


[Install]
WantedBy=default.target
