# dot_config/systemd/user/praise.timer
[Unit]
Description=Praise Mistress

[Timer]
OnCalendar=hourly
Persistent=false

[Install]
WantedBy=timers.target
